# templates/runit.conf
start on runlevel [2345]
stop on runlevel [S016]

respawn
exec /sbin/runsvdir-start
